# pppd.socket
[Unit]
Description=pppd socket for DCNet

[Socket]
ListenStream=7654
Accept=yes

[Install]
WantedBy=sockets.target

