# pppd@.service
[Unit]
Description=DCNet pppd server

[Service]
ExecStart=-/usr/sbin/pppd file /etc/ppp/options.dcnet
StandardInput=socket
