# psmash-pppd.socket
[Unit]
Description=pppd socket for DCNet Power Smash

[Socket]
ListenStream=7656
Accept=yes

[Install]
WantedBy=sockets.target

